# skills/cortex/cortex-server.service
[Unit]
Description=Cortex Server User Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/cortex
ExecStart=%h/cortex/cortex-server
Restart=always
RestartSec=5

# Falls eine .env Datei genutzt wird:
# EnvironmentFile=%h/cortex/.env

# Umgebungsvariablen
Environment="CORTEX_DB_PATH=%h/.openclaw/cortex.db"
Environment="CORTEX_PORT=9123"
Environment="CORTEX_LOG_LEVEL=info"
Environment="CORTEX_RATE_LIMIT=100"
Environment="CORTEX_RATE_LIMIT_WINDOW=1m"
Environment="CORTEX_EMBEDDING_MODEL_PATH=%h/.openclaw/gte-small.gtemodel"

[Install]
WantedBy=default.target